# crontab
0 7 * * * . /app/.env.sh && cd /app && python src/main.py >> /var/log/cron.log 2>&1
